# blog.service
[Unit]
Description=Gunicorn instance to serve Flask blog
After=network.target nginx.service
Requires=nginx.service

[Service]
User=bloguser
Group=bloguser
WorkingDirectory=/home/bloguser/blog
ExecStart=/home/bloguser/blog/venv/bin/gunicorn --workers 4 --bind 0.0.0.0:8000 app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

